# overlay/system/etc/systemd/user/pulseaudio.service.d/pipa-droid-audio.conf
# pipa: PulseAudio module-droid-card dlopens Halium/Android audio libs.
# PULSE_MODULES_DROID_EXTRA_CARD_ARGS: discover merges these into module-droid-card-30
# (sink_format=s24le required for vendor "primary output" mix port on SM8250).
[Service]
Environment=LD_LIBRARY_PATH=/system/lib64:/system/lib:/vendor/lib64:/vendor/lib
Environment=HYBRIS_LD_LIBRARY_PATH=/system/lib64:/vendor/lib64
Environment="PULSE_MODULES_DROID_EXTRA_CARD_ARGS=sink_format=s24le rate=48000 hw_volume=false"
