# overlay/system/etc/NetworkManager/dispatcher.d/90-pipa-networking
#!/bin/bash
# Fix WiFi/Ethernet routing on Xiaomi Pad 6 (pipa).
#
# Root cause: this kernel's FIB rejects the route NetworkManager's l3cfg emits
# for directly-connected networks. NM attaches an all-zeros gateway
# (RTA_GATEWAY=0.0.0.0) to the connected (scope-link) route; the kernel flags
# that as a real AF_INET nexthop and runs fib_check_nh_v4_gw(), whose FIB lookup
# for 0.0.0.0 fails with ENETUNREACH ("Nexthop has invalid gateway"). As a
# result NM installs neither the subnet route nor the default route, so WiFi
# associates and gets a DHCP lease but has no L3 connectivity.
#
# The proper fix is in the kernel (treat 0.0.0.0 RTA_GATEWAY as "no gateway").
# Until a patched kernel is flashed, re-add the routes here via iproute2 (which
# omits the bogus gateway). Routes are added with "proto static" so that
# 02default_route_workaround's "ip route flush proto boot" does not wipe them.

IFACE="$1"
ACTION="$2"

case "$IFACE" in
    wlan*|wlp*|eth*|usb*) ;;
    *) exit 0 ;;
esac

case "$ACTION" in
    up|dhcp4-change|connectivity-change|reapply) ;;
    *) exit 0 ;;
esac

# Address/prefix as configured by NM (e.g. "192.168.8.4/24").
cidr=$(ip -4 -o addr show dev "$IFACE" 2>/dev/null | awk '{print $4; exit}')
[ -n "$cidr" ] || exit 0
ip4="${cidr%%/*}"
pfx="${cidr##*/}"
case "$pfx" in ''|*[!0-9]*) exit 0 ;; esac

# Compute the network address from ip4/pfx (POSIX arithmetic).
IFS=. read -r o1 o2 o3 o4 <<EOF
$ip4
EOF
ip32=$(( (o1 << 24) | (o2 << 16) | (o3 << 8) | o4 ))
mask=$(( (0xffffffff << (32 - pfx)) & 0xffffffff ))
net=$(( ip32 & mask ))
netaddr="$(( (net >> 24) & 255 )).$(( (net >> 16) & 255 )).$(( (net >> 8) & 255 )).$(( net & 255 ))"

# 1) Connected/subnet route (no gateway -> kernel accepts it).
ip route replace "$netaddr/$pfx" dev "$IFACE" proto static scope link src "$ip4" metric 600

# 2) Default route. Prefer the DHCP-provided router, then NM's lease, then .1.
gw="${DHCP4_ROUTERS%% *}"
if [ -z "$gw" ]; then
    conn=$(nmcli -t -f NAME,DEVICE con show --active 2>/dev/null | awk -F: -v d="$IFACE" '$2==d{print $1; exit}')
    [ -n "$conn" ] && gw=$(nmcli -t -f DHCP4.OPTION con show "$conn" 2>/dev/null | sed -n 's/.*routers = //p' | head -1)
fi
[ -z "$gw" ] && gw="${ip4%.*}.1"

if [ -n "$gw" ] && [ "$gw" != "0.0.0.0" ]; then
    ip route replace default via "$gw" dev "$IFACE" proto static metric 600
fi

exit 0
