# overlay/system/etc/profile.d/devicehacks.sh
#!/bin/sh
# Device-specific environment for Xiaomi Pad 6 (pipa)
