# overlay/system/usr/lib/systemd/system/usb-moded.service.d/cleanup-gadget.conf
[Service]
ExecStartPre=-/usr/libexec/usb-gadget-cleanup.sh
