# overlay/system/usr/lib/systemd/system/pipa-flashlight-trigger.service
[Unit]
Description=Set flashlight LED trigger for torch control
Before=graphical-session.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/sh -c 'chmod 0666 /sys/class/leds/flashlight/brightness 2>/dev/null; echo torch0_trigger > /sys/class/leds/flashlight/trigger || true'

[Install]
WantedBy=sysinit.target
