# overlay/system/usr/lib/systemd/system/pipa-boot-wlan.path
[Path]
# /dev/wlan appears after the Android side brings up cnss; re-run when it shows up.
PathExists=/dev/wlan
PathChanged=/dev/wlan

[Install]
WantedBy=multi-user.target
